# smartcomm-dashboard.service
[Unit]
Description=SmartComm Diagnostic Dashboard
After=network-online.target irserver.service argononed.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/smartcomm-dashboard
ExecStart=/usr/bin/python3 /opt/smartcomm-dashboard/dashboard.py
Restart=on-failure
RestartSec=5
StateDirectory=smartcomm-dashboard
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
